# oucher.conf
# oucher

description	"oucher service"

start on filesystem

script
exec /usr/local/bin/oucher
end script
