# oucher.conf
# oucher

description	"oucher service"

start on (startup and started udev)

script
exec /usr/local/bin/oucher
end script
